# install/OS specific/Docker/supervisord.conf
[supervisord]
nodaemon=true

[inet_http_server]
port=*:9001

[program:sshd]
command=/usr/sbin/sshd -D